# systemd/ipfix-enricher.service
[Unit]
Description=IPFIX AS Enricher Service
After=network.target csf.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/ipfix-enricher/ipfix_enricher.py
Restart=always
RestartSec=5
User=root
StandardOutput=journal
StandardError=journal
Nice=-10
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
